# debian/eos-tech-support.postinst
#DEBHELPER#

# Disable the /boot bind mount on dist-upgrade
systemctl disable boot.mount
